# config/autostart/chrome.desktop
[Desktop Entry]
Name=chrome
TryExec=google-chrome
Exec=google-chrome
Terminal=false
Type=Application
StartupNotify=false
